# conf/systemd.service
[Unit]
Description=ZeroTier network controller user interface
After=network.target

[Service]
Type=simple
User=ztncui
Group=ztncui
Environment="__YNH_NODE_LOAD_PATH__"
WorkingDirectory=__FINAL_PATH__/src/
ExecStart=__YNH_NPM__ start

[Install]
WantedBy=multi-user.target
